# contrib/ha-storage/services/chkha.service
[Unit]
Description=check NFS, SMB and DRBD status

[Service]
ExecStart=/bin/bash /etc/drbdha/chkha.sh