# contrib/api-stress-test/roles/virtual-kubelet/install/files/virtual-kubelet.service
[Unit]
Description=virtual kubelet daemon

[Service]
ExecStart=/bin/bash /usr/bin/virtual-kubelet.sh
RestartSec=5s
Restart=on-failure

[Install]
WantedBy=multi-user.target